# data/iio-sensor-proxy.service.in
[Unit]
Description=IIO Sensor Proxy service

[Service]
Type=simple
ExecStart=@sbindir@/iio-sensor-proxy
